# ansible/roles/proxmox/files/80-shutdown-on-battery.rules
KERNEL=="AC", SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/sbin/shutdown -h now"
